# misc/etc/nginx/templates/gunicorn.socket.template
[Unit]
Description=Gunicorn socket

[Socket]
ListenStream=${VIRTUAL_ENV}/bin/gunicorn.sock
SocketUser=${USER}
SocketGroup=nginx
SocketMode=660

[Install]
WantedBy=sockets.target
